# codedeploy/upstart/codedeploy.conf
# upstart configuration

description "golang server"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec /home/ec2-user/gocode/bin/golang_web