# codedeploy/upstart/golang_web.conf
# upstart configuration

description "Golang Web App"
author "Hoohawa"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
exec sudo /home/ubuntu/go/bin/golang_web